# silo_prediction_addon/rootfs/etc/services.d/silo_prediction/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Silo Prediction service stopped"
