# silo_prediction_addon/rootfs/etc/services.d/silo_prediction/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Silo Prediction service..."

exec python3 /app/prediction_service.py
