# silo_prediction_addon/rootfs/etc/services.d/silo_prediction/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Silo Prediction service v2.0..."

exec python3 /app/app.py
